# python_mongodbrs/templates/app/newapp.ini
[program:newapp]
command=python /opt/newapp.py
redirect_stderr=true
priority=10
